# fuzzy_agent_backend/backend/fuzzy-api.service
[Unit]
Description=Fuzzy Logic Voice Tutor API
After=network.target

[Service]
User=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/fuzzy_backend
Environment="PATH=/home/YOUR_USERNAME/fuzzy_backend/venv/bin"
Environment="PORT=5000"
ExecStart=/home/YOUR_USERNAME/fuzzy_backend/venv/bin/gunicorn --bind 0.0.0.0:5000 --workers 2 --timeout 120 app:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
